# systemd/icecast-metadata-updater.service
[Unit]
Description=Aktualizacja metadanych Icecast na podstawie pogody
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/icecast-metadata-updater
ExecStart=%h/icecast-metadata-updater/start_updater.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
